# root/etc/xdg/autostart/docker-browser-autostart.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Default browser (autostart)
Comment=Opens the default browser when the session starts. Set BROWSER_AUTOSTART=false to skip. Log: Desktop/browser-autostart.log under /config.
Exec=/usr/local/bin/docker-browser-autostart.sh
Icon=web-browser
Terminal=false
Categories=Network;WebBrowser;
Hidden=false
X-GNOME-Autostart-enabled=true
StartupNotify=false
